# bot/bot/supervisor/crackdb.conf
[program:crackdb]
command = python %(ENV_APP)s/helpers/crackdb.py
redirect_stderr = true
stdout_logfile = %(ENV_APP)s/supervisor/crackdb.log
autorestart = true