# bot/bot/supervisor/tcbbot.conf
[program:tcbbot]
command = python %(ENV_APP)s/tcbbot.py
redirect_stderr = true
stdout_logfile = %(ENV_APP)s/supervisor/tcbbot.log
autostart=true